# tools/vocab-sync/ops/com.openepcis.vocab-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd template for the weekly vocab-sync upstream re-alignment. TEMPLATE: edit the path below to
  your checkout, then install with:

    cp com.openepcis.vocab-sync.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.openepcis.vocab-sync.plist
    # remove with: launchctl unload ~/Library/LaunchAgents/com.openepcis.vocab-sync.plist

  Runs Mondays at 09:00. The wrapper (sync.sh) is fail-soft: if LM Studio is not running it logs and
  exits cleanly. Output goes to the tool's .cache/sync.log (and the StandardOut/Error paths below).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openepcis.vocab-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/sven/Documents/projects/openepcis/openepcis-dpp-ready/tools/vocab-sync/ops/sync.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/sven/Documents/projects/openepcis/openepcis-dpp-ready/tools/vocab-sync/.cache/sync.launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sven/Documents/projects/openepcis/openepcis-dpp-ready/tools/vocab-sync/.cache/sync.launchd.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
